# infrastructure/packer/files/restic-backup.service
[Unit]
Description=Restic Backup to S3
Documentation=https://restic.net/
OnFailure=backup-notify@%n.service

[Service]
Type=oneshot
User=root
Group=root

# Environment file with Restic configuration
EnvironmentFile=/etc/restic/env

# Execute backup script
ExecStart=/opt/zmanim/backup.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=restic-backup

# Timeout after 2 hours
TimeoutStartSec=7200
