# infrastructure/packer/files/zmanim-db-init.service
[Unit]
Description=Zmanim Database Initialization
Documentation=https://github.com/jcom-dev/zmanim
After=postgresql.service zmanim-firstboot.service
Requires=postgresql.service zmanim-firstboot.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Run as root (scripts use sudo -u postgres internally)
User=root
Group=root

# Step 1: Create user (idempotent SQL)
ExecStart=/usr/bin/sudo -u postgres /usr/bin/psql -f /opt/zmanim/init-db.sql

# Step 2: Create database if not exists (shell script, idempotent)
ExecStartPost=/opt/zmanim/create-db.sh

# Step 3: Enable PostGIS extension (idempotent)
ExecStartPost=/usr/bin/sudo -u postgres /usr/bin/psql -d zmanim -c "CREATE EXTENSION IF NOT EXISTS postgis;"

[Install]
WantedBy=multi-user.target
